# net-mssql/rootfs/etc/s6-overlay/s6-rc.d/mssql-server/finish
#!/command/with-contenv bash

echo "-正在停止 MSSQL 服务..."