# net-mssql/rootfs/etc/s6-overlay/s6-rc.d/mssql-init/up
#!/command/execlineb

/etc/s6-overlay/s6-rc.d/mssql-init/mssql-init.sh