# banban-docker/rootfs/etc/s6-overlay/s6-rc.d/banban/finish
#!/command/with-contenv sh
echo "-正在停止 斑斑低代码 服务..."
