# net-mssql/rootfs/etc/s6-overlay/s6-rc.d/dotnet-app/finish
#!/command/with-contenv bash

echo "-正在停止 .NET 应用程序..."
