# composerize/rootfs/etc/s6-overlay/s6-rc.d/darkhttpd/finish
#!/command/with-contenv sh

echo "+正在停止 darkhttpd 服务..."