# kms/rootfs/etc/s6-overlay/s6-rc.d/darkhttpd/run
#!/command/with-contenv sh

if [ "$WEB" = "true" ]; then
  echo "+正在启动 darkhttpd 服务..."
  exec s6-setuidgid http darkhttpd /vlmcsd/www --port 8080 >/dev/null 2>&1
else
  s6-svc -Od .
  exit 0
fi